# tourboxlinux.service
[Unit]
Description=TourBoxLinux

[Service]
Type=simple
ExecStart=@DESTDIR@/tourboxlinux

[Install]
WantedBy=default.target
