# tourboxlinux.service
[Unit]
Description=TourBoxLinux

[Service]
Type=simple
ExecStart=@PREFIX@/bin/tourboxlinux

[Install]
WantedBy=default.target
